# update/update.service
[Unit]
Description=Atlas updater

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update.sh
